# wf-mon.service
[Unit]
Description=Workflow Monitoring Service
After=network.target docker.service

[Service]
User=root
ExecStart=/usr/bin/env bash /usr/local/bin/run_monitoring.sh
Restart=no
Environment="HOME=/root" "USER=root" "PATH=/usr/local/bin:/usr/bin:/bin"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
